# HW3/linux/sftp_watchd.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/sbin/sftp_watchd

[Install]
WantedBy=default.target