# HW3/linux/sftp_watchd.service
[Unit]
Description=Sftp watchd

[Service]
ExecStart=/usr/sbin/sftp_watchd
Restart=on-failure
RestartSec=1s

[Install]
WantedBy=default.target